# lesson7/logscript.timer
[Unit]
Description=Script sheduller

[Timer]
OnUnitActiveSec=30
Unit=logscript.service

[Install]
WantedBy=multi-user.target
